# docker/php-fpm/supervisord.conf
[supervisord]
nodaemon=true

[program:php7-fpm]
command = php-fpm7.1 -R
user = root
autostart = true